# playbooks/roles/web_interface/templates/vpn-web.service
[Unit]
Description=OpenVPN Web Interface
After=network.target
Requires=network.target

[Service]
Type=simple
User=www-data
Group=www-data
EnvironmentFile=/opt/gcp-ovpn-portal/.env
ExecStart=/usr/local/bin/ovpn-portal
Restart=always
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target