# roles/nginx/templates/runtime-config.service.j2
[Unit]
Description=Configure runtime settings
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/configure-runtime
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target